# systemd/argentina-legal-data-mcp-dataset.timer
[Unit]
Description=Weekly rebuild of the InfoLEG offline dataset

[Timer]
OnCalendar=Sun 04:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
